# repo-mirror/DOD/cron.d/ol7_repo_mirroring
12 0,6,12,18 * * * root cd /root/repo_mirror/ol7 && docker compose up -d

