# repo-mirror/DOD/cron.d/ol8_repo_mirroring
22 0,6,12,18 * * * root cd /root/repo_mirror/ol8 && docker compose up -d

